# meta-wistron-net/recipes-phosphor/fans/phosphor-fan/obmc-host-shutdown@.target
[Unit]
Description=Chassis%i (Soft Power Off)
After=multi-user.target
Requires=obmc-poweroff.service
