# meta-wistron-net/meta-p3box/recipes-wistron-net/fan-fail/fan-fail/obmc-fan-fail@.service
[Unit]
Description=FAN fail operation

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/fan-fail %i
