# meta-wistron-net/recipes-wistron-net/fan-hotswap/fan-hotswap/obmc-fan-hotswap@.service
[Unit]
Description=FAN tray hotswap operation

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/fan-hotswap %i
