# meta-wistron-net/recipes-wistron-net/factory-reset/factory-reset/obmc-factory-reset@.service
[Unit]
Description=OBMC factory reset operation

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/factory-reset
