# meta-wistron-net/recipes-wistron-net/psu-hotswap/psu-hotswap/obmc-psu-hotswap@.service
[Unit]
Description=PSU hotswap operation

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/psu-hotswap %i
